# internal/netconf/testdata/firewall-controller.service
# This file was auto generated for machine: 'e0ab02d2-27cd-5a5e-8efc-080ba80cf258' by app version .
# Do not edit.
[Unit]
Description=Firewall controller - configures the firewall based on k8s resources
After=network.target

[Service]
LimitMEMLOCK=infinity
Environment=KUBECONFIG=/etc/firewall-controller/.kubeconfig
ExecStart=/bin/ip vrf exec vrf104009 /usr/local/bin/firewall-controller --service-ip 10.0.16.2 --private-vrf 3981
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
